# assets/rc.local
#!/bin/bash
/data/home/victron-mqtt-bridge > /dev/null 2>/dev/null &