# assets/rc.local
#!/bin/bash
sleep 20 && /data/victron-mqtt-bridge.sh > /data/victron-mqtt-bridge.log 2>&1 &